# ansible/roles/cndeploy-cardano-node/templates/cnode-tu-restart.service.j2
[Unit]
Description=Cardano Node - Topology Updater - restart cnode.service for topology update

[Service]
Type=oneshot
WorkingDirectory=/opt/cardano/cnode/scripts
ExecStart=/bin/bash -c "/bin/systemctl try-restart cnode.service 2>/dev/null || /usr/bin/systemctl try-restart cnode.service 2>/dev/null"
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cnode-tu-restart
